# init.d/gitlab-irc
#!/bin/sh

# Autogenerated by foreman

### BEGIN INIT INFO
# Provides:             app-irc
# Required-Start:       $syslog
# Required-Stop:        $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    app-irc
# Description:          app-irc
### END INIT INFO

NAME=gitlab-irc
DESC="GitLab IRC Client"

set -e
. /lib/lsb/init-functions

case "${1}" in
	start)
		log_daemon_msg "Starting ${DESC}" "${NAME}"
		if (cd /home/git/gitlab-irc/ && nohup foreman start > /home/git/gitlab-irc/irc.log 2>&1 &); then
                        log_end_msg 0
                else
                        log_end_msg 1
                fi
	;;
	stop)
		log_daemon_msg "Stopping ${DESC}" "${NAME}"
		if (kill $(ps aux | grep -E 'lib/irc.rb|lib/web.rb|foreman: master' | grep -v grep | awk '{print $2}')); then
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	restart)
		${0} stop
		sleep 5
		${0} start
	;;
	*)
		echo "Usage: ${0} {start|stop|restart}" >&2
	;;
esac

exit 0
